# deploy/pi-netzero-pitail.service
[Unit]
Description=pi-netzero Wi-Fi recon server (Pi-Tail)
After=network.target
# No USB-gadget dependency: on Pi-Tail the phone link + DHCP are already
# provided by the image. We only run the app.

[Service]
Type=simple
WorkingDirectory=/opt/pi-netzero
# Flip wlan0 itself into monitor mode (single interface owns the radio — Pi-Tail's
# mon0-alongside-wlan0 makes the two fight over the channel and capture nothing).
# Requires the wlan0 station stanza disabled in /boot/firmware/interfaces so the
# chip is never driven as a station (that association is what crashes the firmware).
Environment=PI_NETZERO_IFACE=wlan0
Environment=PI_NETZERO_RELEASE_RADIO=0
Environment=PI_NETZERO_PORT=8080
ExecStart=/opt/pi-netzero/.venv/bin/python /opt/pi-netzero/main.py
Restart=on-failure
RestartSec=3
# Runs as root: needs raw sockets for scapy and to run mon0up.

[Install]
WantedBy=multi-user.target
